# supervisor.conf
[supervisord]
nodaemon=true

[program:email_verify]
user=root
command=npm run-script start
directory=/app
redirect_stderr=true
autorestart=true
stdout_logfile=/app/log/stdout.log
